# Linux/Automatic Gameplay Recording/.config/systemd/user/gpu-screen-recorder.service
[Unit]
Description=Records the display

[Service]
ExecStart=%h/.local/bin/startRecord

[Install]
WantedBy=default.target
